# s6/services/web/run
#!/command/execlineb -P
cd /opt/pingora-manager/web
/usr/local/bin/bun run start
