# s6/services/web/run
#!/command/with-contenv /bin/sh
cd /opt/pingora-manager/web
/usr/local/bin/bun run init-db.mjs
exec /usr/local/bin/bun run start
