# s6/services/pingora/run
#!/command/execlineb -P
/opt/pingora-manager/proxy
